# autopass.service
[Unit]
Description=Socket activated autopass

[Service]
ExecStart=/home/jokke/autopass.cr/bin/autopass
